# crontab
# News aggregation every 30 minutes
# Fields: minute hour day-of-month month day-of-week command
*/30 * * * * /app/.venv/bin/python /app/backend/main.py --config /app/config.yml
